# scripts/algomlb-sync.timer
[Unit]
Description=Daily schedule for AlgoMLB Data Ingestion
Requires=algomlb-sync.service

[Timer]
# Run every day at 01:00 UTC (07:00 PM MDT / 06:00 PM MST)
# Generates predictions for tomorrow's games before the 8 PM MT night preview card.
OnCalendar=*-*-* 01:00:00
# Spread the load to avoid hitting APIs at the exact same second
RandomizedDelaySec=600
# Keep the schedule persistent across reboots
Persistent=true

[Install]
WantedBy=timers.target
