# scripts/fadegoblin-degen.timer
[Unit]
Description=Timer for FadeGoblin Daily Degen Post
RefuseManualStart=no
RefuseManualStop=no

[Timer]
# Run every day at 16:00 UTC = 10:00 AM MDT / 9:00 AM MST
# Mid-morning post: random parlay for the morning scroll / afternoon degens
OnCalendar=*-*-* 16:00:00
Unit=fadegoblin-degen.service
Persistent=true

[Install]
WantedBy=timers.target
