# scripts/fadegoblin-degen.timer
[Unit]
Description=Timer for FadeGoblin Daily Degen Post
RefuseManualStart=no
RefuseManualStop=no

[Timer]
# Run every day at 18:00 UTC = 12:00 PM MDT / 11:00 AM MST
# Noon post: random parlay for the lunch crowd / afternoon degens
OnCalendar=*-*-* 18:00:00
Unit=fadegoblin-degen.service
Persistent=true

[Install]
WantedBy=timers.target
