# scripts/fadegoblin-recap.service
[Unit]
Description=FadeGoblin Nightly Recap Post (Twitter/X and Bluesky)
After=network.target

[Service]
Type=oneshot
WorkingDirectory=/home/opc/AlgoMLB
Environment="PATH=/home/opc/.local/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
ExecStart=/usr/bin/xvfb-run --auto-servernum --server-args="-screen 0 1280x800x24" env PYTHONPATH=fadegoblin_playwright/src /home/opc/AlgoMLB/.venv/bin/python /home/opc/AlgoMLB/fadegoblin_playwright/src/fadegoblin/main.py --mode recap

[Install]
WantedBy=default.target
